# remote-agent/linux-installer/actionfi-agent.service
[Unit]
Description=ActionFi Remote Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/bin
ExecStart=/usr/local/bin/actionfi-agent
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
